# ops/systemd/rpow-healthcheck.service
[Unit]
Description=rpow watchdog (probe /health, restart on persistent fails)
After=rpow-server.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rpow-healthcheck
